# tuxedo-control-center.install
post_install() {
cat <<-NOTE
To use the Tuxedo Control Center enable tccd and tccd-sleep with the
following commands:

systemctl enable --now tccd.service
systemctl enable --now tccd-sleep.service
NOTE
}
